# docker/supervisor/conf.d/queue-worker.conf
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=sh -c "if [ -f /var/www/html/artisan ]; then php /var/www/html/artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=5 --max-time=3600; else echo 'backend belum di-install (lihat STEP 11)'; sleep 60; fi"
directory=/var/www/html
numprocs=2
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=30
stopsignal=SIGTERM
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
